# Task 1/copy-logs.rules
ACTION=="add", KERNEL=="sd[a-z]*", ENV{ID_BUS}=="usb", SYMLINK+="added_usb", RUN+="/usr/bin/copy-logs-to-usb.sh"