# Task 1/copy-logs.rules
ACTION=="add", KERNEL=="sd[a-z]*", RUN+="/usr/bin/copy-logs-to-usb.sh"